# packaging/systemd/file-mover.service
# ============================================================================
# Background File Mover — systemd unit (PLACEHOLDER)
# ============================================================================
# This is a minimal placeholder. The production unit — with the dedicated
# `mover` service account, RuntimeDirectory/StateDirectory/LogsDirectory,
# RequiresMountsFor for the NFS mounts, and the full hardening block
# (NoNewPrivileges, ProtectSystem=strict, ProtectKernel*, etc.) — is authored
# in Milestone 8. See docs/ROADMAP.md (M8) and docs/CAPTURE.md §7.
#
# The first release uses Type=simple (the app logs when it reaches READY);
# Type=notify is deferred.
# ============================================================================

[Unit]
Description=Background File Mover (durable transfer coordinator)
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/file-mover/venv/bin/python -m file_mover --config /etc/file-mover/file-mover.ini service run
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
